# zypper-dup-download-only.timer
[Unit]
Description=Download available rpm package updates daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target